# scripts/webhook-server.service
[Unit]
Description=GitHub Webhook Server for Booking System Auto-Deployment
Documentation=https://github.com/Jonboyweb/booking-system-only-nextjs
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=door50a-br
Group=door50a-br
WorkingDirectory=/home/door50a-br/htdocs/br.door50a.co.uk

# Environment variables
Environment="NODE_ENV=production"
Environment="WEBHOOK_PORT=9001"
EnvironmentFile=/home/door50a-br/htdocs/br.door50a.co.uk/.env

# Start command
ExecStart=/usr/bin/node /home/door50a-br/htdocs/br.door50a.co.uk/scripts/webhook-server.js

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/door50a-br/htdocs/br.door50a.co.uk/logs
ReadWritePaths=/home/door50a-br/htdocs/br.door50a.co.uk/backups
ReadWritePaths=/tmp

# Resource limits
LimitNOFILE=65536
LimitNPROC=512
MemoryLimit=500M
CPUQuota=50%

# Logging
StandardOutput=append:/home/door50a-br/htdocs/br.door50a.co.uk/logs/webhook-systemd.log
StandardError=append:/home/door50a-br/htdocs/br.door50a.co.uk/logs/webhook-systemd-error.log
SyslogIdentifier=booking-webhook

[Install]
WantedBy=multi-user.target